# scripts/com.tunnelsessions.dailytext.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tunnelsessions.dailytext</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>/Users/kevindrivas/Desktop/projects/tunnel-sessions/scripts/TunnelDailySync.app</string>
    </array>

    <!-- Run at 11:55 AM every day (sync + send takes ~1-2 min) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
